# ubuntu/set-brightness-fp2.conf
start on started lightdm

exec echo 228 > /sys/class/leds/lcd-backlight/brightness
